# project/deploy/supervisor.conf
[group:%(proj_name)s]
programs=uwsgi_%(proj_name)s

[program:uwsgi_%(proj_name)s]
command=%(venv_path)s/env/bin/uwsgi --socket %(proj_name)s.sock --module project.wsgi
directory=%(proj_path)s
user=%(user)s
autostart=true
autorestart=true
redirect_stderr=true
environment=LANG="%(locale)s",LC_ALL="%(locale)s",LC_LANG="%(locale)s"
